# cron.conf
[program:cron]
command=/usr/bin/run-cron.sh
autostart=true
autorestart=false