# cron.conf
[program:cron]
command=/usr/bin/run-cron.sh
autostart=false
autorestart=false