# cron.conf
[program:cron]
command=cron
autostart=true
autorestart=false