# nginx-server.conf
[program:nginxserver]
command=/usr/bin/run-nginx.sh
autostart=true
autorestart=true